# supervisor_services.conf
[program:ircd]
user=ircd
priority=50
autostart=true
autorestart=true
command=ophion -foreground

